# Resources/com.ebowwa.battery-cap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Reference template for the LaunchDaemon that re-applies the battery cap
  on boot. NOT bundled into the binary — PersistenceInstaller generates
  this dynamically with the correct binary path substituted in. Kept here
  for documentation and for hand-install if needed.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ebowwa.battery-cap</string>

    <!-- Path to the BatteryCap binary on the target machine. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/BatteryCap.app/Contents/MacOS/BatteryCap</string>
        <string>--boot-apply</string>
    </array>

    <!-- Re-apply the cap as soon as the daemon loads at boot. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Defense-in-depth: re-apply hourly. Logs to /Library/Logs/BatteryCap.log
         so we can decide later whether this is needed. See PRD OQ1. -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- Don't restart if it exits — short-lived process, runs and exits. -->
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
